# .config/autostart/pidgin-autostart.desktop
[Desktop Entry]
Type=Application
Name=Pidgin Autostart
Comment=Start Pidgin
Exec=pidgin
OnlyShowIn=LXDE
