# .config/autostart/thunderbird-autostart.desktop
[Desktop Entry]
Type=Application
Name=Thunderbird Autostart
Comment=Thunderbird
Exec=thunderbird
OnlyShowIn=LXDE
